# recipes-scripts/scripts/opal6-scripts/initsystem.service
[Unit]
Description=Opal-6 system init

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/InitSystem.sh

[Install]
WantedBy=multi-user.target
