# recipes-scripts/scripts/opal6-scripts/qtstartup.service
[Unit]
Description=QT Demo App

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/qtstartup.sh

[Install]
WantedBy=multi-user.target
