# recipes-scripts/scripts/opal6-scripts/startup.service
[Unit]
Description=Opal-6 startup

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/Startup.sh

[Install]
WantedBy=sysinit.target
